# hi-play-sms-gateway.conf
[program:hi-play-sms-gateway]
command=php /www/GatewayWorker/start.php start -d
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/www/GatewayWorker/tmp/workerman.log